# supervisord.conf
[supervisord]

[program:worker]
command=python3.5 worker.py
process_name=%(program_name)s-%(process_num)s
numprocs=5
directory=.
stopsignal=TERM
autostart=true
autorestart=true
